# configs/ansible/roles/archival-node/templates/systemd/ch-supply.timer.j2
# systemd timer: daily ch-supply token_supply refresh.
#
# Rendered by Ansible. Companion to ch-supply.service. Keeps the lake-derived
# per-token supply (token_supply) fresh so SEP-41 token supply on /v1/assets
# tracks the network instead of freezing at whatever ledger ch-supply last ran.
#
# Time slot: 06:00 UTC + jitter, after compute-completeness (05:30). The seed is
# chunked + memory-guarded so it never collides with the morning oneshots'
# ClickHouse load. Persistent=true: the wrapper's chunked seed absorbs a
# post-outage backlog safely.

[Unit]
Description=Daily ch-supply token_supply refresh (SEP-41 lake-derived supply)
Documentation=https://github.com/StellarIndex/stellar-index/blob/main/docs/architecture/clickhouse-supply-from-ch.md

[Timer]
OnCalendar=*-*-* 06:00:00 UTC
RandomizedDelaySec=300
Persistent=true
Unit=ch-supply.service

[Install]
WantedBy=timers.target
